# docker/s6/worker1/run
#!/bin/bash

/app/worker \
        --slynk-interface "0.0.0.0" \
        --slynk-port "4005" \
        | tee -a /app/logs/worker-1.log \
        | jsail
